# log_analyzer.conf
description     "Test log analyzer service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/test/devartops/log_analyzer.py
